# tests/fixtures/lamp-rhel/usr/lib/systemd/system/httpd.service
[Unit]
Description=The Apache HTTP Server

[Service]
ExecStart=/usr/sbin/httpd
Type=notify

[Install]
WantedBy=multi-user.target
