# tests/fixtures/lamp-rhel/usr/lib/systemd/system/mariadb.service
[Unit]
Description=MariaDB database server

[Service]
ExecStart=/usr/libexec/mysqld
Type=notify

[Install]
WantedBy=multi-user.target
