# tests/fixtures/lamp-rhel/usr/lib/systemd/system/nginx.service
[Unit]
Description=The nginx HTTP and reverse proxy server

[Service]
ExecStart=/usr/sbin/nginx
Type=forking

[Install]
WantedBy=multi-user.target
