# tests/fixtures/java-tomcat/usr/lib/systemd/system/postgresql.service
[Unit]
Description=PostgreSQL
[Service]
ExecStart=/usr/lib/postgresql/14/bin/postgres
[Install]
WantedBy=multi-user.target
