# tests/fixtures/java-tomcat/usr/lib/systemd/system/tomcat.service
[Unit]
Description=Apache Tomcat
[Service]
ExecStart=/opt/tomcat/bin/catalina.sh run
[Install]
WantedBy=multi-user.target
