# tests/fixtures/lamp-debian/usr/lib/systemd/system/apache2.service
[Unit]
Description=The Apache HTTP Server

[Service]
ExecStart=/usr/sbin/apache2
Type=forking

[Install]
WantedBy=multi-user.target
