# tests/fixtures/lamp-debian/usr/lib/systemd/system/postgresql.service
[Unit]
Description=PostgreSQL RDBMS

[Service]
ExecStart=/usr/lib/postgresql/14/bin/postgres
Type=notify

[Install]
WantedBy=multi-user.target
